# systemd/steamlink.service
[Unit]
Description = Steam Link

[Service]
Type = simple
ExecStart = /usr/bin/steamlink
